# .config/systemd/user/log_ip.timer
[Unit]
Description=Run log_ip daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
