# .config/systemd/user/log_ip.service
[Unit]
Description=Log public IP address

[Service]
Type=oneshot
ExecStart=/home/lucho/.local/bin/log_ip
